# provision/roles/zabbix/templates/nasne-status@.timer.j2
[Unit]
Description=Collect nasne status for %I

[Timer]
OnCalendar={{ zabbix_nasne_status_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
